# face-rekon/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:qdrant]
command=/usr/local/bin/qdrant --config-path /app/qdrant-config.yaml
directory=/app
autostart=true
autorestart=false
stdout_logfile=/var/log/qdrant.log
stderr_logfile=/var/log/qdrant-error.log
priority=100
startsecs=5
startretries=1

[program:face-rekon]
command=python3 /app/scripts/app.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/face-rekon.log
stderr_logfile=/var/log/face-rekon-error.log
priority=200
startsecs=10
environment=PYTHONPATH="/app/scripts",QDRANT_HOST="localhost",QDRANT_PORT="6333",FACE_REKON_USE_QDRANT="true"
